# inference-runtime/runit/policy_daemon/run
#!/bin/bash

exec /usr/sbin/policy_daemon
